# deploy/polymarket-bot.service.example
[Unit]
Description=Polymarket Hedge Bot 24/7
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=polymarket
Group=polymarket
WorkingDirectory=/opt/polymarket-agent
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/polymarket-agent/.venv/bin/python -m polymarket_hedge_bot.bot_runner --live-polymarket --stake 200 --live-pages 5 --live-limit 100 --live-min-liquidity 500 --live-orderbook --interval 15 --http-timeout 4 --max-workers 10 --max-loss 200 --min-decision WATCH --min-score 30 --min-edge 0.03 --min-positive-probability 0.50 --min-hours-to-deadline 3 --max-days-to-deadline 45 --min-no-price 0.05 --max-no-price 0.97 --min-net-upside 0 --min-reward-risk 0.10 --radar-max-days-to-deadline 45 --cooldown-min 30 --heartbeat-min 60
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
