# deploy/polymarket-bot.service.example
[Unit]
Description=Polymarket Hedge Bot 24/7
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/polymarket-agent
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/polymarket-agent/.venv/bin/python -m polymarket_hedge_bot.bot_runner --telegram-only --position-monitor --position-monitor-min 1
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
